# backend/ag-prediction.service
[Unit]
Description=Agricultural Dashboard — Price Prediction API
After=network.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/Agricultural_Data_Analysis
Environment="PATH=/home/ec2-user/Agricultural_Data_Analysis/backend/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/home/ec2-user/Agricultural_Data_Analysis/.env
ExecStart=/home/ec2-user/Agricultural_Data_Analysis/backend/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
